# misc/mbp11,1/etc/local.d/bootmisc.start
#!/bin/sh

echo XHC1 > /proc/acpi/wakeup
echo LID0 > /proc/acpi/wakeup
